# systemd/lcp-frontend.service
[Unit]
Description=Readium LCP Server - Test frontend server
#After=network.target

[Service]
#Type=simple
# Another Type option: forking
#User=nanodano
WorkingDirectory="$GOPATH"
ExecStart="$GOPATH/bin/frontend"
Restart=on-failure
# Other Restart options: or always, on-abort, etc

[Install]
WantedBy=multi-user.target
